# config/launchagents/com.aos.memory.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aos.memory</string>
    <key>WorkingDirectory</key>
    <string>__HOME__/aos/core/services/memory</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/.aos/services/memory/.venv/bin/python</string>
        <string>__HOME__/aos/core/services/memory/main.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/.aos/logs/memory.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.aos/logs/memory.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:__HOME__/.bun/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
