# config/launchagents/com.aos.comms-bus.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aos.comms-bus</string>
    <key>WorkingDirectory</key>
    <string>__HOME__/aos</string>
    <key>ProgramArguments</key>
    <array>
        <string>__HOME__/aos/core/bin/internal/comms-bus</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>__HOME__/.aos/logs/comms_bus.out.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/.aos/logs/comms_bus.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:__HOME__/.bun/bin:__HOME__/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>COMMS_POLL_INTERVAL</key>
        <string>300</string>
        <key>COMMS_BUS_PORT</key>
        <string>4099</string>
    </dict>
</dict>
</plist>
